# docker/conf/supervisor-zookeeper.conf
[program:zookeeper]
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/zookeeper/zookeeper.log
directory=/
command=/usr/share/zookeeper/bin/zkServer.sh start-foreground

