# docker/conf/supervisor-ssh.conf
[program:ssh]
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/sshd.log
directory=/
command=/usr/sbin/sshd -De
