# docker/conf/supervisor-bookkeeper.conf
[program:bookkeeper]
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/bookkeeper/bookkeeper.log
directory=/opt/bookkeeper
command=/opt/bookkeeper/bin/bookkeeper bookie

